# board/r36s/rootfs_overlay_systemd/usr/lib/systemd/system/my_init_script.service
[Unit]
Description=My Init Script

[Service]
ExecStart=/usr/local/bin/my_init_script.sh

[Install]
WantedBy=multi-user.target