# board/r36s/rootfs_overlay_systemd/usr/lib/systemd/system/reload-module.service
[Unit]
Description=Reload module after suspend
After=suspend.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/modprobe -r dwc2
ExecStartPost=/usr/sbin/modprobe dwc2

[Install]
WantedBy=suspend.target
